# rootfs/etc/s6-overlay/s6-rc.d/dxcontainer-oneshot--print-license/up
/etc/s6-overlay/s6-rc.d/dxcontainer-oneshot--print-license/run